# packaging/pgcow_0.0.1~alpha1/etc/profile.d/pgcow.sh
#!/usr/bin/env bash

set -e

export PATH="/opt/pgcow/bin:$PATH"
